# supervisor/legion_autostart.conf
[group:legion]
programs = setting,canbus,routing,prediction,planning,control,legion_bridge

[program:setting]
directory=/home/beacon/legion_framework/scripts
command=/bin/bash /home/beacon/legion_framework/scripts/handle_autostart.sh setting
autostart=true
autorestart=true
startsecs=5
startretries=1
stopasgroup=false
redirect_stderr=true
stdout_logfile=/home/beacon/legion_framework/supervisor_logs/setting.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
user=root
priority=1001
environment=
    HOME="/home/beacon",
    USER="beacon"

[program:discovery_server]
directory=/home/beacon/legion_framework/scripts
command=/bin/bash /home/beacon/legion_framework/scripts/handle_autostart.sh discovery_server
autostart=false
autorestart=true
startsecs=5
startretries=3
stopasgroup=true
redirect_stderr=true
stdout_logfile=/home/beacon/legion_framework/supervisor_logs/discovery_server.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
user=beacon
priority=1000
environment=
    HOME="/home/beacon",
    USER="beacon",
    ROS_DOMAIN_ID="0",
    DDS_IMPLEMENTATION="cyclonedds",
    RMW_IMPLEMENTATION="rmw_cyclonedds_cpp",
    RCL_ASSERT_RMW_ID_MATCHES="rmw_cyclonedds_cpp",
    CYCLONEDDS_URI="file:///home/beacon/legion_framework/config/ros_cyclonedds.xml",
    ROS_LOG_DIR="/home/beacon/.ros/log",
    CYCLONEDDS_URI="file:///home/beacon/legion_framework/config/ros_cyclonedds.xml"

[program:zenoh_bridge]
directory=/home/beacon/legion_framework/scripts
command=/bin/bash /home/beacon/legion_framework/scripts/handle_autostart.sh zenoh_bridge
autostart=false
autorestart=true
startsecs=5
startretries=3
stopasgroup=true
redirect_stderr=true
stdout_logfile=/home/beacon/legion_framework/supervisor_logs/zenoh_bridge.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
user=beacon
priority=1002
environment=
    HOME="/home/beacon",
    USER="beacon",
    ROS_DOMAIN_ID="0",
    DDS_IMPLEMENTATION="cyclonedds",
    RMW_IMPLEMENTATION="rmw_cyclonedds_cpp",
    RCL_ASSERT_RMW_ID_MATCHES="rmw_cyclonedds_cpp",
    CYCLONEDDS_URI="file:///home/beacon/legion_framework/config/ros_cyclonedds.xml",
    ROS_LOG_DIR="/home/beacon/.ros/log",
    CYCLONEDDS_URI="file:///home/beacon/legion_framework/config/ros_cyclonedds.xml"

[program:canbus]
directory=/home/beacon/legion_framework/scripts
command=/bin/bash /home/beacon/legion_framework/scripts/handle_autostart.sh canbus
autostart=true
autorestart=true
startsecs=10
startretries=3
stopasgroup=true
redirect_stderr=true
stdout_logfile=/home/beacon/legion_framework/supervisor_logs/canbus.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
user=beacon
priority=999
environment=
    HOME="/home/beacon",
    USER="beacon",
    ROS_DOMAIN_ID="0",
    DDS_IMPLEMENTATION="cyclonedds",
    MY_IP="192.168.1.101",
    MY_INTERFACE="eth0",
    ROS_LOG_DIR="/home/beacon/.ros/log",
    CYCLONEDDS_URI="file:///home/beacon/legion_framework/config/ros_cyclonedds.xml"

[program:routing]
directory=/home/beacon/legion_framework/scripts
command=/bin/bash /home/beacon/legion_framework/scripts/handle_autostart.sh routing
autostart=true
autorestart=true
startsecs=10
startretries=3
stopasgroup=true
redirect_stderr=true
stdout_logfile=/home/beacon/legion_framework/supervisor_logs/routing.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
user=beacon
priority=999
environment=
    HOME="/home/beacon",
    USER="beacon",
    ROS_DOMAIN_ID="0",
    DDS_IMPLEMENTATION="cyclonedds",
    MY_IP="192.168.1.101",
    MY_INTERFACE="eth0",
    ROS_LOG_DIR="/home/beacon/.ros/log",
    CYCLONEDDS_URI="file:///home/beacon/legion_framework/config/ros_cyclonedds.xml"

[program:prediction]
directory=/home/beacon/legion_framework/scripts
command=/bin/bash /home/beacon/legion_framework/scripts/handle_autostart.sh prediction
autostart=true
autorestart=true
startsecs=10
startretries=3
stopasgroup=true
redirect_stderr=true
stdout_logfile=/home/beacon/legion_framework/supervisor_logs/prediction.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
user=beacon
priority=999
environment=
    HOME="/home/beacon",
    USER="beacon",
    ROS_DOMAIN_ID="0",
    DDS_IMPLEMENTATION="cyclonedds",
    MY_IP="192.168.1.101",
    MY_INTERFACE="eth0",
    ROS_LOG_DIR="/home/beacon/.ros/log",
    CYCLONEDDS_URI="file:///home/beacon/legion_framework/config/ros_cyclonedds.xml"

[program:planning]
directory=/home/beacon/legion_framework/scripts
command=/bin/bash /home/beacon/legion_framework/scripts/handle_autostart.sh planning
autostart=true
autorestart=true
startsecs=10
startretries=3
stopasgroup=true
redirect_stderr=true
stdout_logfile=/home/beacon/legion_framework/supervisor_logs/planning.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
user=beacon
priority=999
environment=
    HOME="/home/beacon",
    USER="beacon",
    ROS_DOMAIN_ID="0",
    DDS_IMPLEMENTATION="cyclonedds",
    MY_IP="192.168.1.101",
    MY_INTERFACE="eth0",
    ROS_LOG_DIR="/home/beacon/.ros/log",
    CYCLONEDDS_URI="file:///home/beacon/legion_framework/config/ros_cyclonedds.xml"

[program:control]
directory=/home/beacon/legion_framework/scripts
command=/bin/bash /home/beacon/legion_framework/scripts/handle_autostart.sh control
autostart=true
autorestart=true
startsecs=10
startretries=3
stopasgroup=true
redirect_stderr=true
stdout_logfile=/home/beacon/legion_framework/supervisor_logs/control.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
user=beacon
priority=999
environment=
    HOME="/home/beacon",
    USER="beacon",
    ROS_DOMAIN_ID="0",
    DDS_IMPLEMENTATION="cyclonedds",
    MY_IP="192.168.1.101",
    MY_INTERFACE="eth0",
    ROS_LOG_DIR="/home/beacon/.ros/log",
    CYCLONEDDS_URI="file:///home/beacon/legion_framework/config/ros_cyclonedds.xml"

[program:legion_bridge]
directory=/home/beacon/legion_framework/scripts
command=/bin/bash /home/beacon/legion_framework/scripts/handle_autostart.sh legion_bridge
autostart=true
autorestart=true
startsecs=10
startretries=3
stopasgroup=true
redirect_stderr=true
stdout_logfile=/home/beacon/legion_framework/supervisor_logs/legion_bridge.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
user=beacon
priority=999
environment=
    HOME="/home/beacon",
    USER="beacon",
    ROS_DOMAIN_ID="0",
    DDS_IMPLEMENTATION="cyclonedds",
    MY_IP="192.168.1.101",
    MY_INTERFACE="eth0",
    ROS_LOG_DIR="/home/beacon/.ros/log",
    CYCLONEDDS_URI="file:///home/beacon/legion_framework/config/ros_cyclonedds.xml"

[program:wxapp_server]
directory=/home/beacon/legion_framework/scripts
command=/bin/bash /home/beacon/legion_framework/scripts/handle_autostart.sh wxapp_server
autostart=true
autorestart=true
startsecs=10
startretries=3
stopasgroup=true
redirect_stderr=true
stdout_logfile=/home/beacon/legion_framework/supervisor_logs/wxapp.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
user=beacon
priority=999
environment=
    HOME="/home/beacon",
    USER="beacon",
    ROS_DOMAIN_ID="0",
    DDS_IMPLEMENTATION="cyclonedds",
    MY_IP="192.168.1.101",
    MY_INTERFACE="eth0",
    ROS_LOG_DIR="/home/beacon/.ros/log",
    CYCLONEDDS_URI="file:///home/beacon/legion_framework/config/ros_cyclonedds.xml"
